# deploy/systemd/azul.service
[Unit]
Description=Tabletopia Azul Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/tabletopia/azul/packages/server
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=ALLOWED_ORIGIN=https://game.yuxinlai.com,http://game.yuxinlai.com
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
